# processes.ini
[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:crond]
command=/usr/sbin/crond -n -x proc
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
